# scripts/com.santisplayground.restart.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.santisplayground.restart</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/sh</string>
      <string>-lc</string>
      <string>REPO_DIR="$HOME/santisplayground_repo" /bin/sh "$HOME/santisplayground_repo/scripts/restart.sh"</string>
    </array>

    <!-- Restart every day at 04:00 local time -->
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>4</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/santisplayground-restart.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/santisplayground-restart.err</string>
  </dict>
</plist>

